# systemd/goodwe-master-coordinator-user.service
[Unit]
Description=GoodWe Master Coordinator Service (User)
Documentation=https://github.com/rafalmachnik/goodwe-dynamic-price-optimiser
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
WorkingDirectory=/home/rmachnik/sources/goodwe-dynamic-price-optimiser
ExecStart=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/venv/bin/python /home/rmachnik/sources/goodwe-dynamic-price-optimiser/src/master_coordinator.py --non-interactive
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goodwe-master-coordinator

# Environment variables
Environment=PYTHONPATH=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/src
Environment=PYTHONUNBUFFERED=1
Environment=HOME=/home/rmachnik

# User service - no security restrictions that require root
# These are commented out for user services:
# NoNewPrivileges=true
# PrivateTmp=true
# ProtectSystem=strict
# ProtectHome=true
# ReadWritePaths=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/logs
# ReadWritePaths=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/out
# ReadWritePaths=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/config
# IPAddressAllow=localhost
# IPAddressAllow=10.0.0.0/8
# IPAddressAllow=172.16.0.0/12
# IPAddressAllow=192.168.0.0/16
# LimitNOFILE=65536
# MemoryMax=512M
# CPUQuota=50%

[Install]
WantedBy=default.target