# systemd/goodwe-vps-tunnel.service
[Unit]
Description=GoodWe Master Coordinator VPS SSH Tunnel
After=goodwe-master-coordinator.service
BindsTo=goodwe-master-coordinator.service

[Service]
Type=simple
User=rmachnik
Group=rmachnik
WorkingDirectory=/home/rmachnik/sources/goodwe-dynamic-price-optimiser
# Reverse SSH Tunnel: Expose Pi's local 8080 to VPS local 20358
# Options:
# -o StrictHostKeyChecking=no: Don't fail if host key is not in known_hosts
# -o ExitOnForwardFailure=yes: Fail if the port is already in use on the VPS
# -o ServerAliveInterval=60: Keep connection alive
# Cleanup: Kill any previous tunnel process on the VPS using port 20358
# Cleanup: Kill any previous tunnel process on the VPS using port 20358 (always succeed)
ExecStartPre=/usr/bin/ssh -p 10358 -o StrictHostKeyChecking=no -i /home/rmachnik/.ssh/id_goodwe_tunnel root@andrzej358.mikrus.xyz "fuser -k 20358/tcp 2>/dev/null || kill \$(lsof -t -i:20358) 2>/dev/null || true"
# Reverse SSH Tunnel: Expose Pi's local 8080 to VPS 0.0.0.0:20358
# 0.0.0.0:20358 ensures the Mikrus public proxy can connect to the internal tunnel endpoint
ExecStart=/usr/bin/ssh -N -R 0.0.0.0:20358:localhost:8080 -p 10358 -o StrictHostKeyChecking=no -o ExitOnForwardFailure=yes -o ServerAliveInterval=60 -i /home/rmachnik/.ssh/id_goodwe_tunnel root@andrzej358.mikrus.xyz
# Maintenance
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goodwe-vps-tunnel

[Install]
WantedBy=multi-user.target
