# edge/systemd/goodwe-cloud-reporter.service
[Unit]
Description=GoodWe Cloud Reporter
After=network.target goodwe-master-coordinator.service

[Service]
Type=simple
User=rmachnik
WorkingDirectory=/home/rmachnik/sources/goodwe-dynamic-price-optimiser
Environment=PYTHONPATH=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/edge:/home/rmachnik/sources/goodwe-dynamic-price-optimiser/edge/src
EnvironmentFile=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/.env
ExecStart=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/venv/bin/python edge/cloud_reporter.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
