# edge/systemd/goodwe-master-coordinator.service
[Unit]
Description=GoodWe Master Coordinator Service
Documentation=https://github.com/rafalmachnik/goodwe-dynamic-price-optimiser
After=network.target network-online.target
Wants=network-online.target
StartLimitIntervalSec=0

[Service]
Type=simple
User=rmachnik
Group=rmachnik
WorkingDirectory=/home/rmachnik/sources/goodwe-dynamic-price-optimiser
Environment=PYTHONPATH=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/edge/src
Environment=PYTHONUNBUFFERED=1
EnvironmentFile=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/.env
ExecStart=/home/rmachnik/sources/goodwe-dynamic-price-optimiser/venv/bin/python edge/src/master_coordinator.py --non-interactive
ExecStop=/bin/kill -TERM $MAINPID
KillMode=mixed
KillSignal=TERM
Restart=always
RestartSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=goodwe-master-coordinator

[Install]
WantedBy=multi-user.target
